# postgres_xl/supervisord_master.conf
[supervisord]
nodaemon=true
user=root

[program:gtm]
user=postgres
command=/usr/local/pgsql/bin/gtm -D /usr/local/pgsql/data/data_gtm

[program:master]
user=postgres
command=/usr/local/pgsql/bin/postgres --coordinator -c pooler_port=40100 -D /usr/local/pgsql/data/data_coord1